# deploy/systemd/garmin-backup.service
# OPS-02 · nightly online backup of garmin.db (SQLite VACUUM INTO + rotation).
# Install alongside the bot/uvicorn units (see CLAUDE.md TODO). Triggered by the
# matching .timer — a oneshot, not a long-running service.
#
#   sudo cp deploy/systemd/garmin-backup.* /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now garmin-backup.timer
#
# Adjust WorkingDirectory / User / --rsync-dest for the Pi. The rsync destination is
# the OFF-SD copy — without it a corrupt SD card takes the "backups" with it.
[Unit]
Description=Garmin bot: nightly SQLite backup
After=network-online.target

[Service]
Type=oneshot
User=sergiwez
WorkingDirectory=/home/sergiwez/projects/garmin-coach
# --dir backups: rotated copies on the SD card. --rsync-dest: the OFF-SD copy on a
# USB stick mounted at /mnt/backup (see /etc/fstab, UUID FA41-C949) — without it a
# corrupt SD card takes the "backups" with it.
ExecStart=/home/sergiwez/projects/garmin-coach/venv/bin/python -m scripts.backup_db --dir backups --rsync-dest /mnt/backup/garmin/
