# deploy/garmin-web.service
[Unit]
Description=Garmin Bot — FastAPI web
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/garmin_bot
EnvironmentFile=/home/pi/garmin_bot/.env
ExecStartPre=/home/pi/garmin_bot/venv/bin/python -m alembic upgrade head
ExecStart=/home/pi/garmin_bot/venv/bin/python -m uvicorn app.main:create_app \
    --factory --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
