# deploy/garmin-admin-bot.service
[Unit]
Description=Garmin Admin Bot — hidden /deploy + /test_* commands (owner-only)
After=network.target

[Service]
Type=simple
User=sergiwez
WorkingDirectory=/home/sergiwez/projects/garmin-coach
EnvironmentFile=/home/sergiwez/projects/garmin-coach/.env
ExecStart=/home/sergiwez/projects/garmin-coach/venv/bin/python -m bot.admin_main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
